# liquidcfg.service
[Unit]
Description=AIO Startup Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/liquidcfg.sh

[Install]
WantedBy=default.target
